# deploy/com.kaiju.daily.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
    "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<!--
  com.kaiju.daily.plist — launchd agent for the Kaiju daily lifecycle.

  OPERATOR SETUP: before installing this file, edit the three REPLACE placeholders:

  (1) ProgramArguments path: replace /REPLACE/WITH/ABS/PATH/legendary-kaiju
      with the absolute path to your repo clone, e.g.:
        /Users/yourname/projects/legendary-kaiju/deploy/run_daily.sh

  (2) StartCalendarInterval: 06:30 local is the default. Align with NOAA data
      availability per docs/superpowers/notes/noaa-forecast-contract.md.
      NBM 06Z is typically ready by ~06:00 ET; Kalshi NYC markets open ~06:00 ET.
      Adjust forward if NBM data consistently arrives late.

  (3) StandardOutPath / StandardErrorPath: /tmp/kaiju.daily.log is ephemeral.
      Consider ~/Library/Logs/kaiju.daily.log for persistence across reboots.

  SECRETS: add an EnvironmentVariables dict with your KALSHI_KEY_ID,
  KALSHI_PRIVATE_KEY (PEM string), KAIJU_MODE, and KAIJU_CITY values.
  Never commit secret values; set them directly in the installed plist.

  INSTALL:
    cp deploy/com.kaiju.daily.plist ~/Library/LaunchAgents/
    # edit the installed copy to set secrets and absolute paths
    launchctl load ~/Library/LaunchAgents/com.kaiju.daily.plist

  UNLOAD:
    launchctl unload ~/Library/LaunchAgents/com.kaiju.daily.plist
-->
<dict>

  <key>Label</key>
  <string>com.kaiju.daily</string>

  <!-- REPLACE: set the absolute path to run_daily.sh on this machine.
       Example: /Users/yourname/projects/legendary-kaiju/deploy/run_daily.sh -->
  <key>ProgramArguments</key>
  <array>
    <string>/REPLACE/WITH/ABS/PATH/legendary-kaiju/deploy/run_daily.sh</string>
  </array>

  <!-- Daily at 06:30 local time. REPLACE Hour/Minute to match your NOAA timing. -->
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>6</integer>
    <key>Minute</key>
    <integer>30</integer>
  </dict>

  <!-- REPLACE: set to a writable absolute log path. -->
  <key>StandardOutPath</key>
  <string>/tmp/kaiju.daily.log</string>

  <key>StandardErrorPath</key>
  <string>/tmp/kaiju.daily.log</string>

  <!-- Do NOT run immediately on launchctl load; wait for the scheduled time. -->
  <key>RunAtLoad</key>
  <false/>

  <!--
    SECRETS + PATH — launchd does NOT inherit the user's shell environment.
    You MUST set every value here; a missing PATH means `uv` is not found
    and the job silently fails. Never commit real secret values.

    REPLACE each placeholder before installing:
      PATH             - include the directory where `uv` is installed on this machine.
      KALSHI_KEY_ID    - your Kalshi API key ID.
      KALSHI_PRIVATE_KEY - your Kalshi RSA private key as a single-line PEM string
                           (no literal BEGIN/END header lines; use the raw base64 only,
                            or a path to the key file, whichever your auth code reads).
      KAIJU_MODE       - shadow-paper (safe default) or live (requires KAIJU_LIVE_ARM_TOKEN).
      KAIJU_CITIES     - comma-separated station list read by kaiju.Settings (e.g. KNYC).
      KAIJU_CITY       - single station code read by run_daily.sh for the station flag.
  -->
  <key>EnvironmentVariables</key>
  <dict>
    <!-- REPLACE: colon-separated PATH that includes the directory containing `uv`. -->
    <key>PATH</key>
    <string>/REPLACE/HOME/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>

    <!-- REPLACE: Kalshi API credentials. -->
    <key>KALSHI_KEY_ID</key>
    <string>REPLACE_KALSHI_KEY_ID</string>

    <!-- REPLACE: RSA private key PEM string (single-line or file path per your auth impl). -->
    <key>KALSHI_PRIVATE_KEY</key>
    <string>REPLACE_PEM_STRING</string>

    <!-- REPLACE: shadow-paper is safe; set to live only with KAIJU_LIVE_ARM_TOKEN also set. -->
    <key>KAIJU_MODE</key>
    <string>shadow-paper</string>

    <!-- REPLACE: comma-separated station list used by kaiju.Settings (KAIJU_CITIES env var). -->
    <key>KAIJU_CITIES</key>
    <string>KNYC</string>

    <!-- REPLACE: single station code used by run_daily.sh (station flag) (KAIJU_CITY env var). -->
    <key>KAIJU_CITY</key>
    <string>NYC</string>
  </dict>

</dict>
</plist>
